# deploy/supervisord_web.conf
[supervisord]

[supervisorctl]

[program:web-app]
command = python manage.py runserver 80